# dockerComposeZammad.service
[Unit]
Description=Zammad service with docker compose
Requires=docker.service
After=docker.service

[Service]
Type=simple
RemainAfterExit=yes
WorkingDirectory=/opt/zammadDocker
ExecStart=docker compose -f /opt/zammadDocker/docker-compose.yml up -d --remove-orphans
ExecStop=docker compose down

[Install]
WantedBy=multi-user.target
